# ATYappyBot.service
# for Ubuntu 16.04 or highest
[Unit]
Description=ATYappyBot Demon
After=network.target

[Service]
Type=simple
User=kavlev

Environment=PYTHONPATH=/home/ATYappyBot/
WorkingDirectory=/home/ATYappyBot/
ExecStart=/usr/local/bin/python3.9 /home/ATYappyBot/main.py
Restart=always
RestartSec=60

[Install]
WantedBy=network.target

